# Services/AutostartBrowser.service
[Unit]
Description=Launch Firefox in Kiosk Mode
After=graphical.target

[Service]
ExecStart=/usr/bin/firefox --kiosk https://InfoDisplay:5000
Restart=always

[Install]
WantedBy=graphical.target
