# Services/AutostartBrowser.service
[Unit]
Description=Autostart Browser Service
After=network.target

[Service]
ExecStart=/usr/bin/firefox --kiosk %SERVER_IP%
Restart=always
User=%u
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
