# Services/AutostartServer.service
[Unit]
Description=Autostart Server Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 %INSTALL_DIR%/Server/FlaskServer.py
Restart=always
User=%u

[Install]
WantedBy=default.target
